# deploy/rro.container
# Reason Ready — Podman Quadlet (single node).
#
# A Quadlet is a declarative unit that podman's systemd generator turns into a
# real service. Rootless install (recommended):
#
#   mkdir -p ~/.config/containers/systemd
#   install -m644 deploy/rro.container   ~/.config/containers/systemd/rro.container
#   install -m640 deploy/config.env.example ~/.config/rro/config.env   # then edit
#   systemctl --user daemon-reload
#   systemctl --user start rro
#
# Rootful install: put the file in /etc/containers/systemd/ and use
# `systemctl` (no --user). Verify generation with:
#   /usr/lib/systemd/system-generators/podman-system-generator --dryrun
#
# Build the image first (see deploy/Containerfile). For a model-backed node,
# build with --build-arg FEATURES=candle and mount weights via the Volume
# line below (uncomment), pointing RRO_EMBEDDER/RRO_RERANKER at them.

[Unit]
Description=Reason Ready Objects (RRO) — a2a node
After=network-online.target
Wants=network-online.target

[Container]
Image=localhost/rro:latest
# Pull policy: never for locally-built images; set to newer for a registry.
PullPolicy=never
# The a2a warp-point listener and the ops/metrics port.
PublishPort=7878:7878
PublishPort=9090:9090
# Durable estate lives in a named volume (survives restarts/updates).
Volume=rro-estate.volume:/var/lib/rro:Z
# Model weights, when running candle backends. Mount a host dir of
# .safetensors/tokenizer.json and point RRO_EMBEDDER_WEIGHTS at it.
# Volume=/srv/rro/models:/var/lib/rro/models:ro,Z
EnvironmentFile=%h/.config/rro/config.env
# Clean drain: the daemon commits its RRD baseline on SIGTERM.
Environment=RUST_LOG=info

[Service]
Restart=on-failure
RestartSec=2
TimeoutStopSec=30

[Install]
WantedBy=default.target
